# openclaw.service
# /etc/systemd/system/openclaw.service
# OpenClaw Systemd Service (Native - No Docker)
# ==============================================
#
# YOUR_USER is automatically replaced by install-openclaw.sh during installation.
#
# Installation (handled by install-openclaw.sh):
#   sed "s/YOUR_USER/$(whoami)/g" openclaw.service | sudo tee /etc/systemd/system/openclaw.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable openclaw
#   sudo systemctl start openclaw
#
# Commands:
#   systemctl status openclaw     # Check status
#   sudo systemctl restart openclaw  # Restart
#   journalctl -u openclaw -f     # View logs
#

[Unit]
Description=OpenClaw AI Assistant Gateway
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/opt/openclaw

# Environment
EnvironmentFile=/opt/openclaw/.env.openclaw
Environment=NODE_ENV=production
Environment=OPENCLAW_PREFER_PNPM=1
Environment=OPENCLAW_DISABLE_BONJOUR=1
Environment=PATH=/home/YOUR_USER/.bun/bin:/usr/local/bin:/usr/bin:/bin

# Pre-start cleanup
ExecStartPre=/opt/openclaw/startup-cleanup.sh

# Main process - run OpenClaw gateway directly with Node
ExecStart=/usr/bin/node /opt/openclaw/dist/index.js gateway --allow-unconfigured

# Graceful stop
ExecStop=/bin/kill -SIGTERM $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Timeouts
TimeoutStartSec=30
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/openclaw /var/log/openclaw /home/YOUR_USER/.openclaw

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
